# ansible/roles/lmod/templates/lmod-scache.timer.j2
[Unit]
Description="Update lmod system spider cache."

[Timer]
OnCalendar={{ lmod_timer_oncalender }}
Persistent=true
Unit=lmod-scache.service

[Install]
WantedBy=timers.target
